# os/overlays/etc/inittab
# BusyBox init for BigFred hub
::sysinit:/etc/init.d/rcS
::respawn:-/sbin/getty -L tty1 0 vt100
::respawn:-/sbin/getty -L ttyAMA10 115200 vt100
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
