# os/overlays/etc/crontabs/root
# Hub log rotation (§8.9)
0  3  *   *   *     /usr/sbin/rotate-hub-logs
15 3  *   *   *     /usr/sbin/rotate-hub-logs
